# deploy/clab/dailyhothole-snapshot.timer
[Unit]
Description=Publish the DailyHotHole public snapshot every two hours

[Timer]
OnBootSec=10min
OnCalendar=*-*-* 00/2:05:00
RandomizedDelaySec=5min
Persistent=true
AccuracySec=1min
Unit=dailyhothole-snapshot.service

[Install]
WantedBy=timers.target
