# deploy/clab/dailyhothole-snapshot.service
[Unit]
Description=Generate and publish the DailyHotHole public snapshot
After=network-online.target pkuhole-monitor.service
Wants=network-online.target

[Service]
Type=oneshot
User=pkuhole
Group=pkuhole
WorkingDirectory=/opt/dailyhothole-public/repo
Environment=DAILYHOTHOLE_REPO_DIR=/opt/dailyhothole-public/repo
Environment=DAILYHOTHOLE_SOURCE_URL=http://127.0.0.1:8766/api/state
Environment=DAILYHOTHOLE_SSH_KEY=/opt/dailyhothole-public/ssh/github-deploy-key
Environment=DAILYHOTHOLE_KNOWN_HOSTS=/opt/dailyhothole-public/ssh/known_hosts
Environment=DAILYHOTHOLE_MAX_DAYS=30
Environment=DAILYHOTHOLE_TOP_N=10
Environment=DAILYHOTHOLE_MAX_COMMENTS=500
ExecStart=/opt/dailyhothole-public/publish-snapshot.sh
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/opt/dailyhothole-public/repo
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
UMask=0027
